# agent/os/macos/installer/pkg_build/scripts/postinstall
#!/usr/bin/env bash
# =============================================================================
#  postinstall — AttackLens Agent v2.0.0
#
#  Reads optional env file before installing:
#    echo "ATTACKLENS_MANAGER='<IP_or_URL>'" > /tmp/attacklens_envs
#    echo "ATTACKLENS_AGENT_NAME='<name>'"   >> /tmp/attacklens_envs
#
#  Supported env vars:
#    ATTACKLENS_MANAGER    — manager IP, IP:port, or full http/https URL
#    ATTACKLENS_AGENT_NAME — human-readable agent label (default: ComputerName)
#    ATTACKLENS_TAGS       — comma-separated tags, e.g. "prod,finance"
#    ATTACKLENS_TOKEN      — enrollment token (if manager requires one)
# =============================================================================
set -uo pipefail

readonly INSTALL_DIR="/Library/AttackLens"
readonly BIN_DIR="/Library/AttackLens/bin"
readonly LOG_DIR="/Library/AttackLens/logs"
readonly SECURITY_DIR="/Library/AttackLens/security"
readonly CONFIG_PATH="/Library/AttackLens/agent.toml"
readonly LDIR="/Library/LaunchDaemons"
readonly PKG_DEFAULT_MANAGER="http://YOUR_MANAGER_IP:8080"
readonly VERSION="2.0.0"

log()  { echo "  $1"; }
die()  { echo "  ERROR: $1" >&2; exit 1; }
ok()   { echo "  [OK] $1"; }
warn() { echo "  [WARN] $1"; }

echo ""
echo "  AttackLens Agent v${VERSION} — Installing..."
echo "  =================================================="

# ── Step 1: Read /tmp/attacklens_envs ─────────────────────────────────────────
log "Reading environment..."
ATTACKLENS_MANAGER=""
ATTACKLENS_AGENT_NAME=""
ATTACKLENS_TAGS=""
ATTACKLENS_TOKEN=""

ENV_FILE="/tmp/attacklens_envs"
if [[ -f "$ENV_FILE" ]]; then
    log "  Found: ${ENV_FILE}"
    while IFS= read -r line || [[ -n "$line" ]]; do
        # Skip comments and blank lines
        [[ "$line" =~ ^[[:space:]]*# ]] && continue
        [[ -z "${line//[[:space:]]/}" ]] && continue
        key="${line%%=*}"
        val="${line#*=}"
        # Strip surrounding single or double quotes
        val="${val#\'}" ; val="${val%\'}"
        val="${val#\"}"  ; val="${val%\"}"
        val="${val#\'}" ; val="${val%\'}"   # handle echo "KEY='val'" format
        case "$key" in
            ATTACKLENS_MANAGER)    ATTACKLENS_MANAGER="$val" ;;
            ATTACKLENS_AGENT_NAME) ATTACKLENS_AGENT_NAME="$val" ;;
            ATTACKLENS_TAGS)       ATTACKLENS_TAGS="$val" ;;
            ATTACKLENS_TOKEN)      ATTACKLENS_TOKEN="$val" ;;
        esac
    done < "$ENV_FILE"
    [[ -n "$ATTACKLENS_MANAGER"    ]] && log "  Manager   : $ATTACKLENS_MANAGER"
    [[ -n "$ATTACKLENS_AGENT_NAME" ]] && log "  Agent name: $ATTACKLENS_AGENT_NAME"
    [[ -n "$ATTACKLENS_TAGS"       ]] && log "  Tags      : $ATTACKLENS_TAGS"
else
    log "  No /tmp/attacklens_envs — using package defaults"
fi

# ── Step 2: Normalise manager URL ─────────────────────────────────────────────
# Priority: env file > pkg baked-in default
MANAGER="${ATTACKLENS_MANAGER:-$PKG_DEFAULT_MANAGER}"

# If bare IP or IP:PORT (no scheme), prepend http://
if [[ -n "$MANAGER" ]] && [[ "$MANAGER" != http://* ]] && [[ "$MANAGER" != https://* ]]; then
    MANAGER="http://${MANAGER}"
fi
log "  Effective manager URL: ${MANAGER}"

# ── Step 3: Permissions ───────────────────────────────────────────────────────
log "Setting permissions..."
chown -R root:wheel "${INSTALL_DIR}"
chmod 755 "${INSTALL_DIR}" "${BIN_DIR}"
chmod 644 "${BIN_DIR}/run_agent.py" "${BIN_DIR}/run_watchdog.py"
chmod 755 "${BIN_DIR}/generate_config.sh"
chmod 750 "${LOG_DIR}" "${INSTALL_DIR}/data" "${INSTALL_DIR}/spool"
chmod 700 "${SECURITY_DIR}"
for cli in /usr/local/bin/attacklens-service /usr/local/bin/attacklens; do
    [[ -e "$cli" ]] && { chown root:wheel "$cli"; chmod 755 "$cli"; } || true
done
for f in com.attacklens.agent.plist com.attacklens.watchdog.plist; do
    [[ -f "${LDIR}/${f}" ]] && chown root:wheel "${LDIR}/${f}" && chmod 644 "${LDIR}/${f}" || true
done
ok "Permissions"

# ── Step 4: Find real Mach-O python3 ──────────────────────────────────────────
# macOS 15+ Sequoia/Tahoe: ProgramArguments[0] MUST be a Mach-O binary.
# Ask Python to resolve its own real path (handles all symlink/wrapper chains).
log "Detecting python3..."
PYTHON3=""
for p in /usr/local/bin/python3 /opt/homebrew/bin/python3 /usr/bin/python3; do
    [[ -x "$p" ]] || continue
    real="$("$p" -c "import os,sys; print(os.path.realpath(sys.executable))" 2>/dev/null)" || continue
    [[ -x "$real" ]] || continue
    if file "$real" 2>/dev/null | grep -q "Mach-O"; then
        PYTHON3="$real"; break
    fi
done
# Fallback: accept any python3 (may be a wrapper; will work but may hit exit-78 on macOS 15+)
if [[ -z "$PYTHON3" ]]; then
    for p in /usr/local/bin/python3 /opt/homebrew/bin/python3 /usr/bin/python3; do
        [[ -x "$p" ]] && { PYTHON3="$p"; break; }
    done
fi
[[ -z "$PYTHON3" ]] && die "python3 not found. Install from https://python.org then reinstall."
ok "python3: ${PYTHON3}  ($("${PYTHON3}" --version 2>&1))"

# Warn if Python < 3.11
"${PYTHON3}" -c "import sys; exit(0 if sys.version_info>=(3,11) else 1)" 2>/dev/null ||     warn "Python < 3.11 — upgrade recommended: https://python.org/downloads/"

# ── Step 5: Install Python prerequisites (always, not conditional) ────────────
log "Installing Python prerequisites..."
_pip() {
    local pkg="$1"
    "${PYTHON3}" -m pip install --quiet --upgrade --break-system-packages "${pkg}" 2>/dev/null ||     "${PYTHON3}" -m pip install --quiet --upgrade "${pkg}" 2>/dev/null ||     warn "Could not install ${pkg} — agent may have reduced functionality"
}
# tomllib is stdlib in 3.11+; install tomli as backport for older Python
"${PYTHON3}" -c "import tomllib" 2>/dev/null || _pip tomli
_pip psutil
_pip cryptography
_pip requests
ok "Prerequisites installed"

# ── Step 6: Patch plists with real Mach-O python3 ────────────────────────────
log "Patching LaunchDaemon plists..."
for plist in "${LDIR}/com.attacklens.agent.plist" "${LDIR}/com.attacklens.watchdog.plist"; do
    [[ -f "$plist" ]] || continue
    sed -i '' "s|PYTHON3_PLACEHOLDER|${PYTHON3}|g" "$plist"
    # Also replace any previously-set path that differs from detected Mach-O
    "${PYTHON3}" - "$plist" "${PYTHON3}" <<'PYEOF'
import sys, re
plist, py3 = sys.argv[1], sys.argv[2]
with open(plist) as f: content = f.read()
new = re.sub(
    r'(<key>ProgramArguments</key>\s*<array>\s*<string>)[^<]*(</string>)',
    rf'\g<1>{py3}\g<2>', content, flags=re.DOTALL)
if new != content:
    with open(plist, 'w') as f: f.write(new)
PYEOF
    ok "Patched: $(basename "$plist") → ${PYTHON3}"
done

# ── Step 7: Generate agent.toml ───────────────────────────────────────────────
# Always regenerate when env file provides values; otherwise skip if exists.
SHOULD_GEN=false
[[ ! -f "${CONFIG_PATH}" ]] && SHOULD_GEN=true
[[ -f "${ENV_FILE}" ]]       && SHOULD_GEN=true   # env file present → always apply

if [[ "$SHOULD_GEN" == "true" ]]; then
    log "Generating agent.toml..."
    MANAGER_URL="${MANAGER}"                               AGENT_NAME="${ATTACKLENS_AGENT_NAME:-}"                AGENT_TAGS="${ATTACKLENS_TAGS:-}"                      ENROLL_TOKEN="${ATTACKLENS_TOKEN:-}"                   INSTALL_DIR="${INSTALL_DIR}"                           LOG_DIR="${LOG_DIR}"                                   SECURITY_DIR="${SECURITY_DIR}"                         bash "${BIN_DIR}/generate_config.sh"
    ok "agent.toml written"
else
    log "  Existing agent.toml kept (no env file provided)"
    log "  To reconfigure: sudo attacklens-service set-manager <IP>"
fi

# ── Step 8: Validate plists ───────────────────────────────────────────────────
for plist in "${LDIR}/com.attacklens.agent.plist" "${LDIR}/com.attacklens.watchdog.plist"; do
    plutil -lint "$plist" >/dev/null 2>&1         && ok "Valid plist: $(basename "$plist")"         || warn "Plist syntax error: $(basename "$plist")"
done

# ── Step 9: Start services in SYSTEM domain (never gui/) ─────────────────────
log "Starting LaunchDaemons..."
for LABEL in com.attacklens.agent com.attacklens.watchdog; do
    # Make sure not stuck in wrong domain from previous installs
    for uid in $(dscl . -list /Users UniqueID 2>/dev/null | awk '$2>499{print $2}'); do
        launchctl bootout "gui/${uid}/${LABEL}" 2>/dev/null || true
    done
    launchctl enable "system/${LABEL}" 2>/dev/null || true
done

launchctl bootstrap system "${LDIR}/com.attacklens.agent.plist" 2>/dev/null     && ok "Started: com.attacklens.agent"     || log "  com.attacklens.agent queued — run: sudo attacklens-service start"

launchctl bootstrap system "${LDIR}/com.attacklens.watchdog.plist" 2>/dev/null     && ok "Started: com.attacklens.watchdog"     || log "  com.attacklens.watchdog queued (non-critical)"

# ── Step 10: Verify & summary ─────────────────────────────────────────────────
sleep 3
AGENT_STATE=$(launchctl print "system/com.attacklens.agent" 2>/dev/null | awk '/state =/{print $3}')

echo ""
echo "  =================================================="
echo "  AttackLens Agent v${VERSION} — Installation Complete"
echo "  =================================================="
echo "  Manager : ${MANAGER}"
echo "  Config  : ${CONFIG_PATH}"
echo "  Python  : ${PYTHON3}"
echo "  Agent   : ${AGENT_STATE:-checking...}"
echo ""
echo "  Service commands:"
echo "    attacklens-service status     # check agent state"
echo "    sudo attacklens-service start  # start"
echo "    sudo attacklens-service stop   # stop"
echo "    sudo attacklens-service restart # restart"
echo "    attacklens-service logs        # live log"
echo "    attacklens-service diagnose    # connectivity check"
echo ""
if [[ "$MANAGER" == *"YOUR_MANAGER"* ]]; then
    echo "  NOTE: Set manager IP:"
    echo "    sudo attacklens-service set-manager <IP>"
    echo ""
fi
# Clean up env file (contains potentially sensitive token)
rm -f "${ENV_FILE}" 2>/dev/null && log "  Cleaned up ${ENV_FILE}"
