# agent/os/macos/installer/pkg_build/scripts/postinstall
#!/usr/bin/env bash
# =============================================================================
#  postinstall — AttackLens Agent v2.0.0
#  Runs as root on the target Mac after the pkg payload is laid down.
# =============================================================================
set -uo pipefail

readonly INSTALL_DIR="/Library/AttackLens"
readonly BIN_DIR="/Library/AttackLens/bin"
readonly LOG_DIR="/Library/AttackLens/logs"
readonly SECURITY_DIR="/Library/AttackLens/security"
readonly CONFIG_PATH="/Library/AttackLens/agent.toml"
readonly LDIR="/Library/LaunchDaemons"
readonly DEFAULT_MANAGER_URL="http://34.224.174.38:8080"

log() { echo "  $1"; }
die() { echo "  ERROR: $1" >&2; exit 1; }

# ── 1. Set file ownership & permissions ───────────────────────────────────────
log "Setting permissions..."
chown -R root:wheel "${INSTALL_DIR}"
chmod 755 "${INSTALL_DIR}" "${BIN_DIR}"
chmod 644 "${BIN_DIR}/run_agent.py" "${BIN_DIR}/run_watchdog.py"
chmod 755 "${BIN_DIR}/generate_config.sh"
chmod 750 "${LOG_DIR}" "${INSTALL_DIR}/data" "${INSTALL_DIR}/spool"
chmod 700 "${SECURITY_DIR}"

for cli in /usr/local/bin/attacklens-service /usr/local/bin/attacklens; do
    [[ -e "$cli" ]] && { chown root:wheel "$cli"; chmod 755 "$cli"; } || true
done

for f in com.attacklens.agent.plist com.attacklens.watchdog.plist; do
    chown root:wheel "${LDIR}/${f}" 2>/dev/null || true
    chmod 644        "${LDIR}/${f}" 2>/dev/null || true
done

# ── 2. Detect python3 (MUST be a Mach-O binary, not a shell script) ───────────
log "Detecting python3..."
PYTHON3=""
for p in /usr/local/bin/python3 /opt/homebrew/bin/python3 /usr/bin/python3; do
    if [[ -x "$p" ]]; then
        # Reject shell-script wrappers — macOS 15+ LaunchDaemon requirement
        if file "$p" 2>/dev/null | grep -q "Mach-O"; then
            PYTHON3="$p"; break
        else
            log "  Skipping shell-script python3 at $p"
        fi
    fi
done

if [[ -z "$PYTHON3" ]]; then
    # Last resort: accept any executable python3 (handles wrapper chains)
    for p in /usr/local/bin/python3 /opt/homebrew/bin/python3 /usr/bin/python3; do
        [[ -x "$p" ]] && { PYTHON3="$p"; break; }
    done
fi

[[ -z "$PYTHON3" ]] && die "python3 not found. Install from https://python.org and reinstall."
log "  python3: ${PYTHON3}  ($("${PYTHON3}" --version 2>&1))"

# ── 3. Verify Python version ≥ 3.11 ──────────────────────────────────────────
PY_OK=$("${PYTHON3}" -c "
import sys
ok = (sys.version_info >= (3,11))
print('yes' if ok else 'no')
" 2>/dev/null || echo "no")
if [[ "$PY_OK" != "yes" ]]; then
    log "WARNING: Python < 3.11 detected. Agent requires 3.11+."
    log "         Download: https://python.org/downloads/"
    log "         Continuing install — agent may fail to start until upgraded."
fi

# ── 4. Patch PYTHON3_PLACEHOLDER in LaunchDaemon plists ─────────────────────
log "Patching plists..."
for plist in "${LDIR}/com.attacklens.agent.plist" "${LDIR}/com.attacklens.watchdog.plist"; do
    if [[ -f "$plist" ]]; then
        sed -i '' "s|PYTHON3_PLACEHOLDER|${PYTHON3}|g" "$plist"
        log "  Patched: $(basename "$plist")"
    fi
done

# ── 5. Install Python dependencies ────────────────────────────────────────────
log "Checking Python dependencies..."

# tomllib (stdlib in 3.11+) or tomli (backport for 3.10 and earlier)
"${PYTHON3}" -c "import tomllib" 2>/dev/null || "${PYTHON3}" -c "import tomli"   2>/dev/null || {
    log "  Installing tomli (TOML parser backport)..."
    "${PYTHON3}" -m pip install --quiet --break-system-packages tomli 2>/dev/null ||     "${PYTHON3}" -m pip install --quiet tomli 2>/dev/null ||     log "  WARNING: Could not install tomli — agent may fail to read config on Python < 3.11"
}

# psutil (optional but strongly recommended — falls back to CLI tools without it)
"${PYTHON3}" -c "import psutil" 2>/dev/null || {
    log "  Installing psutil (system metrics collector)..."
    "${PYTHON3}" -m pip install --quiet --break-system-packages psutil 2>/dev/null ||     "${PYTHON3}" -m pip install --quiet psutil 2>/dev/null ||     log "  NOTE: psutil not installed — using CLI fallbacks (lsof, sysctl, etc.)"
}

# cryptography (required for AES-256-GCM payload encryption)
"${PYTHON3}" -c "from cryptography.hazmat.primitives.ciphers.aead import AESGCM" 2>/dev/null || {
    log "  Installing cryptography..."
    "${PYTHON3}" -m pip install --quiet --break-system-packages cryptography 2>/dev/null ||     "${PYTHON3}" -m pip install --quiet cryptography 2>/dev/null ||     log "  WARNING: cryptography not installed — payload encryption disabled"
}

# requests (HTTP transport)
"${PYTHON3}" -c "import requests" 2>/dev/null || {
    log "  Installing requests..."
    "${PYTHON3}" -m pip install --quiet --break-system-packages requests 2>/dev/null ||     "${PYTHON3}" -m pip install --quiet requests 2>/dev/null ||     log "  WARNING: requests not installed — agent cannot send telemetry"
}

# ── 6. Write complete agent.toml (skip if already configured) ─────────────────
if [[ ! -f "${CONFIG_PATH}" ]]; then
    log "Generating agent.toml..."
    MANAGER_URL="${DEFAULT_MANAGER_URL}"     INSTALL_DIR="${INSTALL_DIR}"             LOG_DIR="${LOG_DIR}"                     SECURITY_DIR="${SECURITY_DIR}"           bash "${BIN_DIR}/generate_config.sh"
else
    log "Existing agent.toml found — skipping config generation."
    log "  Run 'sudo attacklens-service update-config' to regenerate."
fi

# ── 7. Validate plist syntax ──────────────────────────────────────────────────
log "Validating plists..."
for plist in "${LDIR}/com.attacklens.agent.plist" "${LDIR}/com.attacklens.watchdog.plist"; do
    if ! plutil -lint "$plist" >/dev/null 2>&1; then
        log "  WARNING: plist syntax error in $(basename $plist) — service may not load"
        plutil -lint "$plist" 2>&1 | sed 's/^/    /' || true
    else
        log "  Valid: $(basename "$plist")"
    fi
done

# ── 8. Load services ──────────────────────────────────────────────────────────
log "Loading LaunchDaemons..."
for LABEL in com.attacklens.agent com.attacklens.watchdog; do
    launchctl enable "system/${LABEL}" 2>/dev/null || true
done

AGENT_STARTED=false
if launchctl bootstrap system "${LDIR}/com.attacklens.agent.plist" 2>/dev/null; then
    AGENT_STARTED=true
    log "  Started: com.attacklens.agent"
elif launchctl load -w "${LDIR}/com.attacklens.agent.plist" 2>/dev/null; then
    AGENT_STARTED=true
    log "  Started: com.attacklens.agent (legacy load)"
else
    log "  NOTE: com.attacklens.agent could not start automatically."
    log "        On first install macOS may require background service approval."
    log "        → System Settings › Privacy & Security › Background Items"
    log "        → Allow items from AttackLens"
    log "        Then run: sudo attacklens-service start"
fi

if launchctl bootstrap system "${LDIR}/com.attacklens.watchdog.plist" 2>/dev/null ||    launchctl load -w "${LDIR}/com.attacklens.watchdog.plist" 2>/dev/null; then
    log "  Started: com.attacklens.watchdog"
else
    log "  NOTE: com.attacklens.watchdog did not start (non-critical — agent runs independently)."
fi

# ── 9. Summary ────────────────────────────────────────────────────────────────
sleep 2
echo ""
echo "  ╔═══════════════════════════════════════════════════════════╗"
echo "  ║  AttackLens Agent v2.0.0 — Installation Complete          ║"
echo "  ╚═══════════════════════════════════════════════════════════╝"
echo ""
echo "  Manager:  ${DEFAULT_MANAGER_URL}"
echo "  Config:   ${CONFIG_PATH}"
echo "  Python:   ${PYTHON3}"
echo ""
echo "  Quick commands:"
echo "    attacklens-service status          — check service status"
echo "    attacklens-service diagnose        — connectivity check"
echo "    attacklens-service logs            — tail live log"
echo ""
if [[ "$DEFAULT_MANAGER_URL" == *"YOUR_MANAGER"* ]]; then
    echo "  ⚠  Manager IP is a placeholder."
    echo "     Set it now:  sudo attacklens-service set-manager <IP>"
    echo ""
fi
echo "  Change manager:"
echo "    sudo attacklens-service set-manager <NEW_IP>"
echo ""
